# board/olimex/imx233_olinuxino_wordclock/skeleton/etc/profile.d/locale.sh
#!/bin/sh
#
# Set LANG
#

export LANG=de_DE.UTF-8
