# board/olimex/imx233_olinuxino_wordclock/skeleton/etc/profile.d/locale.sh
#!/bin/sh
#
# Set LANG and TZ
#

export LANG=de_DE.UTF-8
export TZ=$(cat /etc/TZ)
